# CNV/web-server/src/main/resources/installation/hill-server.service
[Unit]

Description=Initiate java server on port 8000

[Service]
ExecStart=/home/ec2-user/init.sh

[Install]
WantedBy=multi-user.target
